# supervisord.ini
[supervisord]
nodaemon=true

[program:shadowsocks]
command=exec /usr/local/bin/ssserver
    -s localhost -p 8388
    -k %(ENV_SS_KEY)s
    -m %(ENV_SS_METHOD)s
    -t $(ENV_SS_TIMEOUT)d
    %(ENV_SS_ARGS)s

[program:kcptun]
command=exec /usr/local/bin/ktserver
    -t localhost:8388
    -l %(ENV_KT_LISTEN)s
    --key %(ENV_KT_KEY)s
    --mode %(ENV_KT_MODE)s
    --crypt %(ENV_KT_CRYPT)s
    %(ENV_KT_ARGS)s

